# lib/init/kiln-controller.service
[Unit]
Description=kiln-controller

[Service]
ExecStart=/home/kiln-controller/venv/bin/python /home/kiln-controller/kiln-controller.py

[Install]
WantedBy=multi-user.target
